# sshd.conf
[program:sshd]
user=finroc_user
command=/usr/sbin/sshd -D
autorestart=true
stopsignal=INT
